# services/data_collector.service
[Unit]
Description=Daily Sensex Data Collector
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sensex-options-trading-system
Environment=PROJECT_PATH=/home/ubuntu/sensex-options-trading-system
EnvironmentFile=/home/ubuntu/sensex-options-trading-system/.env

# Run data collection at market close (3:30 PM IST)
ExecStartPre=/bin/bash -c 'if [ $(date +%%H%%M) = 1530 ]; then python3 download_weekly_data.py; fi'
ExecStart=/bin/sleep infinity
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
