# services/postback_server.service
[Unit]
Description=Zerodha Postback Server (Venv)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/main_trading
Environment=PATH=/home/ubuntu/main_trading/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/ubuntu/main_trading
Environment=VIRTUAL_ENV=/home/ubuntu/main_trading/venv
ExecStart=/home/ubuntu/main_trading/venv/bin/python3 /home/ubuntu/main_trading/postback_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=postback-server

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/main_trading

# Process limits for free tier
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
