# services/trading_system.service
[Unit]
Description=Sensex Options Trading System (Venv)
After=network-online.target postback_server.service
Wants=network-online.target
Requires=postback_server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/main_trading
Environment=PATH=/home/ubuntu/main_trading/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/ubuntu/main_trading
Environment=VIRTUAL_ENV=/home/ubuntu/main_trading/venv
EnvironmentFile=/home/ubuntu/main_trading/.env

# Dynamic mode reading and execution
ExecStart=/bin/bash -c 'MODE=$(cat /home/ubuntu/main_trading/.trading_mode 2>/dev/null || echo "TEST"); if [ "$MODE" != "DISABLED" ]; then /home/ubuntu/main_trading/venv/bin/python3 /home/ubuntu/main_trading/main.py --mode $(echo $MODE | tr A-Z a-z); else echo "Trading disabled, sleeping..."; sleep infinity; fi'

ExecStop=/bin/bash -c 'echo "DISABLED" > /home/ubuntu/main_trading/.trading_mode'
ExecStopPost=/bin/bash -c 'if [ -f /home/ubuntu/main_trading/.trading_disabled ]; then rm -f /home/ubuntu/main_trading/.trading_disabled; fi'

Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=3

# Resource limits for EC2 free tier
MemoryLimit=512M
CPUQuota=70%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-system

[Install]
WantedBy=multi-user.target
